# systemd/azure-upload.service
[Unit]
Description=Pi Edge Sensing: Azure uploader
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
User=gerrit
WorkingDirectory=/home/gerrit/Projects/pi-sensing
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/gerrit/Projects/pi-sensing/.venv/bin/python /home/gerrit/Projects/pi-sensing/src/uploader.py
Restart=always
RestartSec=10
[Install]
WantedBy=multi-user.target
