# systemd/data-collector.service
[Unit]
Description=Pi Edge Sensing: data collector (UTC minute sampling)
After=network-online.target pigpiod.service
Wants=network-online.target
[Service]
Type=simple
User=gerrit
WorkingDirectory=/home/gerrit/Projects/pi-sensing
Environment=GPIO_BACKENDS=lgpio,rpi
Environment=PULSE_SKIP_PIGPIO=1
Environment=LGPIO_CHIP_PRIORITY=4,0,10,11,12,13
ExecStart=/home/gerrit/Projects/pi-sensing/.venv/bin/python /home/gerrit/Projects/pi-sensing/src/collector.py
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target